# ~s/.config/autostart/kde.sh.desktop
[Desktop Entry]
Exec=/home/l/.dwm/kde.sh
Icon=dialog-scripts
Name=kde.sh
Path=
Type=Application
X-KDE-AutostartScript=true
